# Consul/Server/Supervisord/conf.d/consul-server.conf
[program:consul-server]
command=consul agent -config-file=/etc/consul/server.json
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/1